# IOTstack/custom/owm/crontab
*/5 * * * * /usr/local/bin/python /usr/src/app/app.py
